# deploy/systemd/smartlocker-device-registry.service
[Unit]
Description=Smart Locker Device Registry
After=network-online.target smartlocker-bootstrap.service
Wants=network-online.target
Requires=smartlocker-bootstrap.service

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1
Environment=SMARTLOCKER_CONFIG_JSON=/etc/smartlocker/config.json
Environment=SMARTLOCKER_DEVICE_ID_FILE=/etc/smartlocker/device_id
Environment=SMARTLOCKER_BACKEND_STATE_FILE=/var/lib/smartlocker/backend_device.json
ExecStartPre=/usr/bin/test -x __PROJECT_DIR__/venv/bin/python
ExecStartPre=/usr/bin/test -d /etc/smartlocker
ExecStart=/bin/sh -c 'cd __PROJECT_DIR__ && exec __PROJECT_DIR__/venv/bin/python -m uvicorn app.main:app --app-dir __PROJECT_DIR__ --host 0.0.0.0 --port 8000'
Restart=always
RestartSec=3
StateDirectory=smartlocker
LogsDirectory=smartlocker
StandardOutput=append:/var/log/smartlocker/device-registry.log
StandardError=append:/var/log/smartlocker/device-registry.log

[Install]
WantedBy=multi-user.target
