# app/systemmd/qbox-device.service
[Unit]
Description=Smart Locker Device Runtime (all agents, one shared MQTT connection)
After=network-online.target NetworkManager.service dbus.service bluetooth.service mediamtx.service
Wants=network-online.target NetworkManager.service bluetooth.service mediamtx.service
Conflicts=qbox-telemetry-agent.service qbox-heartbeat-agent.service qbox-control-agent.service qbox-wifi-agent.service qbox-streaming-agent.service wifi-upload-agent.service qbox-streams.service

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=__PROJECT_DIR__
Environment=DEVICE_MODE=production
Environment=DJANGO_API_BASE=https://backend.qbox.sa
Environment=STREAM_AUTO_START=true
Environment=MEDIAMTX_HOST=backend.qbox.sa
Environment=MEDIAMTX_RTSP_PORT=8554
Environment=STREAM_PUBLIC_BASE_URL=https://backend.qbox.sa

ExecStart=__PROJECT_DIR__/.venv/bin/python -m app.main

Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

# WiFi, BLE, camera, and network access required by the combined runtime.
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_FOWNER CAP_CHOWN CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_DAC_OVERRIDE CAP_FOWNER CAP_CHOWN CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_NET_BIND_SERVICE
ProtectSystem=false
ProtectHome=false
PrivateTmp=true
PrivateDevices=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

DeviceAllow=/dev/rfkill rw
DeviceAllow=/dev/ttyAMA0 rw
DeviceAllow=/dev/video0 rw
DeviceAllow=/dev/video1 rw
DeviceAllow=/dev/video2 rw
DeviceAllow=/dev/video3 rw

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
