# app/systemmd/qbox-telemetry-agent.service
[Unit]
Description=Smart Locker Telemetry Agent (standalone/debug)
After=network-online.target
Wants=network-online.target
Conflicts=qbox-device.service

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1

# Standalone service for debugging. In production prefer qbox-device.service
# so telemetry shares the single process-wide MQTT connection with all agents.
ExecStart=__PROJECT_DIR__/.venv/bin/python -m app.agents.telemetry_agent

Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
