# app/systemmd/qbox-streams.service
[Unit]
Description=Smart Locker Streaming Agent (FFmpeg + MQTT)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__

Environment=PYTHONUNBUFFERED=1
Environment=STREAM_AUTO_START=true
Environment=MQTT_HOST=69.62.125.223
Environment=MQTT_PORT=1883
Environment=MEDIAMTX_HOST=backend.qbox.sa
Environment=MEDIAMTX_RTSP_PORT=8554
Environment=STREAM_PUBLIC_BASE_URL=https://backend.qbox.sa
# Uncomment and set if MQTT auth needed
# Environment=MQTT_USERNAME=qbox
# Environment=MQTT_PASSWORD=strongpassword123

ExecStart=__PROJECT_DIR__/.venv/bin/python -m app.streaming_agent.main

Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

# Access to video devices (camera)
DeviceAllow=/dev/video0 rw
DeviceAllow=/dev/video1 rw
DeviceAllow=/dev/video2 rw
DeviceAllow=/dev/video3 rw

# Prevent CPU loops
CPUQuota=80%

# Output
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
