# app/systemmd/wifi-upload-agent.service
[Unit]
Description=Smart Locker WiFi Upload Agent (MQTT + Scan)
After=network-online.target
Wants=network-online.target

# 🔥 IMPORTANT: depend on reconnect service
After=wifi-reconnect.service
Wants=wifi-reconnect.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/smart-locker-hardware

Environment=PYTHONUNBUFFERED=1
Environment=MQTT_HOST=127.0.0.1
Environment=MQTT_PORT=1883

ExecStart=/opt/smart-locker-hardware/.venv/bin/python -m app.hardware_agent.main

Restart=always
RestartSec=10

# prevent CPU loops on network failure
StartLimitIntervalSec=300
StartLimitBurst=10

[Install]
WantedBy=multi-user.target