# app/systemmd/wifi-reconnect.service
[Unit]
Description=Smart Locker WiFi Reconnect Service
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__PROJECT_DIR__/.venv/bin/python -m app.services.reconnect
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
